# pantrypilot_safe/rootfs/etc/services.d/pantrypilot/run
#!/command/with-contenv bash
cd /app
./run.sh
